# install/wtr/wtr.service
[Unit]
Description=Watering Server
After=multi-user.target

[Service]
Type=idle
ExecStart=/wtr/wtr.sh

[Install]
WantedBy=multi-user.target

